# dependence/src/dbus-1.12.20/bus/dbus.socket
[Unit]
Description=D-Bus System Message Bus Socket

[Socket]
ListenStream=/customer/var/run/dbus/system_bus_socket
